# deploy/parser-worker.service
# Example systemd unit file for parser-worker.
#
# Install:
#   sudo install -m 755 ./parser-worker /usr/local/bin/parser-worker
#   sudo install -m 644 deploy/parser-worker.service /etc/systemd/system/
#   sudo install -m 600 .env /etc/parser-worker.env       # adjust path/perms
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now parser-worker.service
#   sudo journalctl -u parser-worker -f                   # follow logs
#
# Set CHECK_INTERVAL_MINUTES > 0 in /etc/parser-worker.env so the binary
# stays alive and polls on its own. With CHECK_INTERVAL_MINUTES=0 systemd
# would just keep restarting the one-shot binary; if that's what you want,
# pair it with a `parser-worker.timer` unit instead.

[Unit]
Description=parser-worker (CSC demo parser → fragg-3.0 stats ingest)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/parser-worker
EnvironmentFile=/etc/parser-worker.env
Restart=on-failure
RestartSec=10s

# Run as a non-privileged dynamic user; systemd creates and tears it down
# automatically. Drop this and supply User=/Group= if you want a fixed UID.
DynamicUser=yes

# Filesystem hardening — the worker only needs to write to /tmp.
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# Modest resource caps. Demos can be ~500 MB and 7z decompression takes
# CPU; bump these if you raise MAX_MATCHES_PER_RUN aggressively.
MemoryMax=2G
TasksMax=128

# Send stdout/stderr to the journal as JSON-summary log lines.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
